# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/html/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:laravel-queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/queue.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/html/storage/logs/queue-error.log
stderr_logfile_maxbytes=10MB
priority=30

[program:laravel-reverb]
command=php /var/www/html/artisan reverb:start --host=0.0.0.0 --port=8080 --no-interaction
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/reverb.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/www/html/storage/logs/reverb-error.log
stderr_logfile_maxbytes=10MB
priority=40
